# hamradio-copilot.service
[Unit]
Description=Hamradio Data Analysis Service
After=network.target

[Service]
User=user
WorkingDirectory=/opt/hamradio-copilot
ExecStart=/usr/bin/python3 /opt/hamradio-copilot/start_copilot.sh
Restart=always

[Install]
WantedBy=multi-user.target
