# hamradio-copilot-web.service
[Unit]
Description=Hamradio Data Analysis Service
After=network.target

[Service]
User=root
Type=simple
WorkingDirectory=/opt/hamradio-copilot
ExecStart=/opt/hamradio-copilot/start_copilot.sh
Restart=always

[Install]
WantedBy=multi-user.target

